# deploy/codeaudit-dashboard.service
[Unit]
Description=CodeAudit Pro — Dashboard Ejecutivo
Documentation=https://github.com/jqime/tech-debt-security-auditor
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/codeauditpro
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 app/dashboard/app.py
Restart=always
RestartSec=5

# Protecciones de seguridad
NoNewPrivileges=yes
ProtectHome=yes
PrivateDevices=yes
ProtectSystem=full

[Install]
WantedBy=multi-user.target
